# buttonreader/buttonreader.service
[Unit]
Description=GPIO Button Reader

[Service]
ExecStart=/usr/bin

[Install]
WantedBy=multi-user.target
